# mywash_admin/supervisord.conf
[supervisord]
nodaemon=false

[program:uwsgi]
command = nohup uwsgi --threads 2 --ini /etc/uwsgi/apps-enabled/dashboard.ini --http-socket 0.0.0.0:9090
autostart=true
autorestart=true
stopsignal=INT